# ignitions/common/files/etc/systemd/system/docker.service.d/10-wait-mount.conf
[Unit]
Requires=var-lib-docker.mount var-lib-kubelet.mount
After=var-lib-docker.mount var-lib-kubelet.mount
