# ignitions/common/files/etc/systemd/system/docker.service.d/90-proxy.conf
[Service]
Environment="HTTP_PROXY=http://localhost:30128" "HTTPS_PROXY=http://localhost:30128"
