# ignitions/common/files/etc/systemd/system/docker.service.d/90-proxy.conf
[Service]
Environment="HTTP_PROXY=http://squid.internet-egress.svc.cluster.local:3128"
Environment="HTTPS_PROXY=http://squid.internet-egress.svc.cluster.local:3128"
