# ignitions/common/files/etc/systemd/system/docker.service.d/90-proxy.conf
[Service]
Environment="HTTP_PROXY=http://{{ index .Spec.IPv4 0 }}:30128"
Environment="HTTPS_PROXY=http://{{ index .Spec.IPv4 0 }}:30128"
