# ignitions/common/files/etc/systemd/system/docker.service.d/90-proxy.conf
[Service]
Environment="HTTP_PROXY=http://127.0.0.1:30128"
Environment="HTTPS_PROXY=http://127.0.0.1:30128"
