# ignitions/common/files/etc/systemd/system/docker.service.d/99-kubernetes-network.conf
[Service]
Environment=DOCKER_OPTS="--bridge=none --iptables=false --ip-masq=false"
