# ignitions/common/files/etc/systemd/system/rkt-gc.timer.d/10-wait-mount.conf
[Unit]
Requires=var-lib-rkt.mount
After=var-lib-rkt.mount
