# ignitions/common/files/etc/systemd/system/rngd.service.d/non_vm.conf
[Unit]
DefaultDependencies=no
ConditionVirtualization=false
