# debian/etc/profile.d/loki-env.sh
export LOKI_ADDR="http://query-frontend.logging.svc:3100/"
export LOKI_ORG_ID="fake"
