# ignitions/roles/cs/systemd/update-neco-operation-cli.service
[Unit]
Description=Update Neco operation CLI commands

[Service]
Type=oneshot
ExecStart=/opt/sbin/update-neco-operation-cli
