# debian/lib/systemd/system/export-unit-status.service
[Unit]
Description=Export statuses of systemd units

[Service]
Type=oneshot
ExecStart=/usr/sbin/export-unit-status
