# debian/lib/systemd/system/rkt-gc.timer
[Unit]
Description=run rkt gc periodically

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
