# debian/lib/systemd/system/rkt-gc.timer
[Unit]
Description=run rkt gc periodically
PartOf=rkt-gc.service

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
