# gcp/public/etc/systemd/system/docker.service.d/99-placemat.conf
[Service]
ExecStartPost=/sbin/iptables -P FORWARD ACCEPT
