# ignitions/common/systemd/chrony-monitor.service
[Unit]
Description=export chrony metrics

[Service]
Type=oneshot
ExecStart=/opt/sbin/chrony-monitor
