# ignitions/common/systemd/docker-prune.timer
[Unit]
Description=run docker prune periodically

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
