# ignitions/common/systemd/neco-dhcp-online.target
[Unit]
Description=Network is Online via DHCP
