# ignitions/common/systemd/neco.target
[Unit]
Description=Wait for routing information
