# ignitions/common/systemd/rkt-gc.service
[Unit]
Description=rkt gc
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/rkt gc
