# ignitions/common/systemd/setup-node-exporter.service
[Unit]
Description=Setup node exporter

[Service]
Type=oneshot
ExecStart=/opt/sbin/setup-node-exporter
RemainAfterExit=yes
