# ignitions/common/systemd/setup-serf.service
[Unit]
Description=Setup serf container
ConditionPathExists=!/etc/serf/serf.json
Wants=neco.target
After=neco.target

[Service]
Type=oneshot
ExecStart=/opt/sbin/setup-serf-conf
RemainAfterExit=yes
