# ignitions/common/systemd/setup-network.service
[Unit]
Description=Setup network
Requires=docker.socket
After=docker.socket

[Service]
Type=oneshot
ExecStart=/opt/sbin/setup-local-network
RemainAfterExit=yes
StandardOutput=journal+console
StandardError=journal+console
