# ignitions/common/systemd/squid-image.service
[Unit]
Description=Load Squid container image
Wants=neco.target
Requires=docker.socket
After=neco.target docker.socket

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/opt/bin/load-image {{ MyURL }}/api/v1/assets/{{ Metadata "squid.img" }} {{ Metadata "squid.ref" }}

[Install]
WantedBy=multi-user.target
