# ignitions/common/systemd/udev-trigger.service
[Unit]
Description=Request udev device events
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/usr/bin/udevadm trigger
