# debian/DEBIAN/prerm
#!/bin/sh -e

SERVICES="neco-updater neco-worker node-exporter sabakan-state-setter neco-rebooter cke cke-localproxy teleport-node trigger-reboot-all-nodes"
TIMERS="docker-prune kill-old-login-sessions backup-cke-etcd export-unit-status trigger-reboot-all-nodes"

prerm() {
    for s in $SERVICES; do
        systemctl stop ${s}.service
        systemctl disable ${s}.service
    done

    for t in $TIMERS; do
        systemctl stop ${t}.timer
        systemctl disable ${t}.timer
    done
}

prerm

exit 0
