# ignitions/common/files/etc/systemd/system/systemd-pstore.service.d/after-local-fs-mount.conf
[Unit]
# systemd-pstore should run after mounting /sys/fs/psotre, which is done in local-fs.target
After=local-fs.target
