# ignitions/common/files/etc/systemd/system/docker.socket.d/99-run-after-time-sync.conf
[Unit]
After=time-sync.target
Requires=time-sync.target
DefaultDependencies=no
