# ignitions/common/files/etc/systemd/system/containerd.service.d/10-wait-mount.conf
[Unit]
Requires=var-lib-containerd.mount
After=var-lib-containerd.mount
