# gcp/public/etc/profile.d/go.sh
export GOPATH=$HOME/go
export GO111MODULE=on
export PATH=$HOME/go/bin:/usr/local/go/bin:/sbin:/usr/sbin:$PATH
